# Personal/.config/systemd/user/docker.service.d/override.conf
[Service]
#Environment="HTTP_PROXY=socks5://10.160.131.136:1080"
#Environment="HTTPS_PROXY=socks5://10.160.131.136:1080"
#Environment="NO_PROXY=localhost,127.0.0.1"
